# scripts/sync-ipad.service
[Unit]
Description=Automatische iPad-Synchronisation bei USB-Verbindung
After=default.target

[Service]
Type=oneshot
ExecStart=/media/Linux/repos/ipad/scripts/sync-ipad.sh